# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

WEBTV_IP=""
DNS_IP="9.9.9.9"

iptables -A INPUT -s 192.168.32.0/24 -p tcp --dport 22 -j REJECT

iptables -t nat -A PREROUTING -d 10.0.0.0/16 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 192.168.102.0/24 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 204.254.74.0/24 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 207.76.180.0/24 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.0/24 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
#iptables -t nat -A PREROUTING -d 210.150.22.0/24 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 210.150.22.35 -s 192.168.32.0/24 -j DNAT --to-destination $DNS_IP
iptables -t nat -A PREROUTING -d 210.150.22.65 -s 192.168.32.0/24 -j DNAT --to-destination $DNS_IP
iptables -t nat -A PREROUTING -d 157.57.221.51 -s 192.168.32.0/24 -j DNAT --to-destination $DNS_IP

## Webone Proxy
#iptables -t nat -A PREROUTING -d 192.168.32.200 -s 192.168.32.0/24 -j DNAT --to-destination 127.0.0.1
#iptables -t nat -A PREROUTING -p tcp -d 192.168.32.200 --dport 8080 -j DNAT --to-destination 127.0.0.1:8080
#iptables -t nat -I OUTPUT -p tcp -d 192.168.32.200  --dport 8080 -j REDIRECT --to-ports 8080
#iptables -t nat -I OUTPUT -p tcp -d 192.168.33.200 --dport 8080 -j DNAT --to-destination 10.0.0.31:8080
#iptables -t nat -I PREROUTING --src 0/0 --dst 192.168.1.1 -p tcp --dport 80 -j REDIRECT --to-ports 8080
#iptables -t nat -I OUTPUT --src 0/0 --dst 192.168.1.1 -p tcp --dport 80 -j REDIRECT --to-ports 8080
iptables -t nat -I PREROUTING --src 0/0 --dst 192.168.31.31 -p tcp --dport 8080 -j REDIRECT --to-ports 8080
iptables -t nat -I OUTPUT --src 0/0 --dst 192.168.31.31 -p tcp --dport 8080 -j REDIRECT --to-ports 8080

iptables -I FORWARD -d 192.168.7.0/24 -s 192.168.32.0/24 -j DROP
#iptables -I FORWARD -d 10.0.0.0/24 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.128/25 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.64/26 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.32/27 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.16/28 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.8/29 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.4/30 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.3/32 -s 192.168.32.0/24 -j DROP

iptables -t nat -A POSTROUTING -s 192.168.32.0/24 -j LOG --log-prefix='[netfilter] '
iptables -t nat -A POSTROUTING -s 192.168.32.0/24 -o end0 -j MASQUERADE

exit 0
