# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

WEBTV_IP=""
DNS_IP="9.9.9.9"

iptables -A INPUT -s 192.168.32.0/24 -p tcp --dport 22 -j REJECT

iptables -t nat -A PREROUTING -d 10.0.0.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.0.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.1.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.1.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.1.3 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.1.4 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 10.0.128.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.1.129 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.128.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.129.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.129.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.130.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.130.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.131.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.131.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.132.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 10.0.132.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 192.168.102.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 192.168.102.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 204.254.74.126 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 207.76.180.1 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 207.76.180.2 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 207.76.180.12 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 207.76.180.13 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 207.76.180.96 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 209.240.194.40 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.41 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.42 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.70 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.71 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.72 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.194.73 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 210.150.22.37 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 210.150.22.58 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 209.240.213.74 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 207.79.35.172 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP
iptables -t nat -A PREROUTING -d 207.79.35.167 -s 192.168.32.0/24 -j DNAT --to-destination $WEBTV_IP

iptables -t nat -A PREROUTING -d 210.150.22.35 -s 192.168.32.0/24 -j DNAT --to-destination $DNS_IP
iptables -t nat -A PREROUTING -d 210.150.22.65 -s 192.168.32.0/24 -j DNAT --to-destination $DNS_IP
iptables -t nat -A PREROUTING -d 157.57.221.51 -s 192.168.32.0/24 -j DNAT --to-destination $DNS_IP

iptables -I FORWARD -d 192.168.7.0/24 -s 192.168.32.0/24 -j DROP
#iptables -I FORWARD -d 10.0.0.0/24 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.128/25 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.64/26 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.32/27 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.16/28 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.8/29 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.4/30 -s 192.168.32.0/24 -j DROP
iptables -I FORWARD -d 10.0.0.3/32 -s 192.168.32.0/24 -j DROP

iptables -t nat -A POSTROUTING -s 192.168.32.0/24 -j LOG --log-prefix='[netfilter] '
iptables -t nat -A POSTROUTING -s 192.168.32.0/24 -o eth0 -j MASQUERADE

exit 0
